# files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:gitit]
user=gitit
command=/usr/bin/gitit -f ./default.conf
directory=/data
